# aaurbs.timer
[Unit]
Description=Run aaurbs every hour

[Timer]
OnCalendar=*:0

[Install]
WantedBy=timers.target
