# aaurbs.service
[Unit]
Description=aaurbs

[Service]
Type=simple
User=aur
Group=aur
ExecStart=/srv/aaurbs/aaurbs.py
