# aaurbs.service
[Unit]
Description=aaurbs

[Service]
Type=simple
ExecStart=/root/git/aaurbs/aaurbs.py
