# aaurbs_webserver.service
[Unit]
Description=aaurbs_webserver

[Service]
Type=simple
ExecStart=/root/git/aaurbs/aaurbs_webserver.py
