# aaurbs.service
[Unit]
Description=aaurbs

[Service]
Type=simple
ExecStart=/aur/updateall.sh
