# aaurbs.service
[Unit]
Description=aaurbs

[Service]
Type=simple
User=aur
Group=aur
EnvironmentFile=/srv/aaurbs/environment.systemd
ExecStart=/srv/aaurbs/aaurbs.py
